# led-terrarium-pi.conf
[program:led-terrarium-pi]
directory=/home/pi/led-terrarium-pi
command=python main.py
autostart=true
autorestart=true
stdout_logfile=/var/log/led-terrarium-pi.log